# doc/openelp.service.in
[Unit]
Description=Open Source EchoLink Proxy
After=network.target nss-lookup.target

[Service]
Type=forking
ConditionPathExists=@FULL_SYSCONF_INSTALL_DIR@/ELProxy.conf
ExecStart=@FULL_BIN_INSTALL_DIR@/openelpd -q -S @FULL_SYSCONF_INSTALL_DIR@/ELProxy.conf

[Install]
WantedBy=multi-user.target
